# php/apache/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:apache2]
autostart=true
command=/usr/sbin/apache2ctl -D FOREGROUND
