# web/python/django/etc/supervisor/conf.d/supervisord.conf
[supervisord]
nodaemon=true

[program:python]
autostart=true
command=python manage.py runserver .
