# old/etc/supervisor/conf.d/mod_php.conf
[include]
files = supervisord.conf

[program:apache2]
autostart=true
command=/usr/bin/pidproxy /var/run/apache2/apache2.pid /usr/local/sbin/apache2
